# examples/fail2ban-digest.cron
0 3 * * *	root	/usr/local/bin/fail2ban_digest maildigest sshd
